# apache2/cron.d/reload
*/5 * * * * root /bin/bash /opt/auto_reload.sh
